# deploy/gunicorn.service
[Unit]
Description=Gunicorn service for Automobile Manufacturing Dashboard
After=network.target

[Service]
User=ec2-user
Group=nginx
WorkingDirectory=/opt/automobile-manufacturing-app
Environment="PATH=/opt/automobile-manufacturing-app/.venv/bin"
ExecStart=/opt/automobile-manufacturing-app/.venv/bin/gunicorn \
    --workers 3 \
    --timeout 120 \
    --bind unix:/run/automobile-manufacturing-app.sock \
    --access-logfile - \
    --error-logfile - \
    wsgi:app
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
