# deploy/gunicorn.service
[Unit]
Description=Gunicorn instance for Automobile Manufacturing Dashboard
After=network.target

[Service]
User=appuser
Group=www-data
WorkingDirectory=/opt/forgepoint
EnvironmentFile=/opt/forgepoint/.env
ExecStart=/opt/forgepoint/venv/bin/gunicorn --workers 3 --bind 127.0.0.1:8001 --access-logfile - --error-logfile - wsgi:app
Restart=always
RestartSec=3

# Hardening (safe defaults for a Flask app under systemd)
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
